# supervisord.conf
[supervisord]
logfile=/dev/stdout 
logfile_maxbytes=0  
loglevel=info
pidfile=/tmp/supervisord.pid
nodaemon=true
user=root

[program:backup]
directory=/root
command=/root/entrypoint.sh

[program:typedb]
directory=/typedb-cluster
command=/typedb-cluster/typedb cluster
